# root/etc/profile.d/build-vars.sh
export PATH=$PATH:/x86_64-pc-linux-musl/bin:/nenuzhnix-tools
export LD_LIBRARY_PATH=/x86_64-pc-linux-musl/lib
export PKG_CONFIG_SYSROOT_DIR=/x86_64-pc-linux-musl
export PKG_CONFIG_LIBDIR=/x86_64-pc-linux-musl/lib/pkgconfig
